# nvidia_fan_prof.service
[Unit]
Description=nVidia Fan Profiling Script
[Service]
ExecStart=/usr/bin/nvidiafan.sh
