# nvidia_fan_prof.timer
[Unit]
Description=nVidia Fan Profiler
[Timer]
OnBootSec=20
OnUnitInactiveSec=4
AccuracySec=500ms
[Install]
WantedBy=timers.target
